# i3de/etc/udev/rules.d/00-keyboards.rules
ACTION=="add", SUBSYSTEM=="input", RUN+="/usr/local/bin/keyboard-mapping.sh"

