# auto-hdmi/95-monitors.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/guy/.local/bin/auto-disp.sh"
